# other-to-delete/supervisord.custom.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=1MB       ; (max main logfile bytes b4 rotation;default 1MB)
logfile_backups=5          ; (num of main logfile rotation backups;default 5)
umask=0000
user=root
;loglevel=debug             

[program:sshd]
command=service ssh start
startsecs=0

[program:watch-log-files]
command=sh /usr/local/bin/super_log_files.sh

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
loglevel=debug

[program:delay]
command=sh /usr/local/bin/delay.sh
priority=7

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
loglevel=debug

[program:post-startup-script]
command=bash /home/dev/startup.sh
autostart=true
autorestart=false
startretries=0

[program:wordpress-security-check]
command=sh /usr/local/bin/wordpress-security-check.sh
autostart=true
autorestart=true
startretries=100

;Kusto telemetry config placeholder
[program:wp-telemetry]
command=bash /usr/local/bin/wordpress-kusto-telemetry.sh
autostart=true
autorestart=true